# monica_crm/rootfs/etc/services.d/apache2/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Monica CRM
# Take down the S6 supervision tree when the apache2 daemon fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "apache2 crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "apache2 stopped, restarting..."