# monica_crm/rootfs/etc/services.d/apache2/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Monica Personal CRM
# Runs the apache server that starts monica
# ==============================================================================

exec httpd -DFOREGROUND
