# monica_crm/rootfs/etc/services.d/apache2/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Monica Personal CRM
# Runs the apache server that starts monica
# ==============================================================================

declare certfile
declare keyfile
declare certpath
declare keypath
declare devserver
declare servername
declare serveradmin

declare enable_access_logs

bashio::log.info "updating apache config with paths to SSL certs"
certfile=$(bashio::config 'certfile')
keyfile=$(bashio::config 'keyfile')
servername=$(bashio::config 'server_name')
serveradmin=$(bashio::config 'server_admin')
devserver=$(bashio::config 'dev_server')
enable_access_logs=$(bashio::config 'enable_access_logs')

if [ "${devserver}" = "true" ]; then
    certpath="/data/${certfile}"
    keypath="/data/${keyfile}"
else
    certpath="/ssl/${certfile}"
    keypath="/ssl/${keyfile}"
fi

if [ "${enable_access_logs}" = "true" ]; then
    bashio::log.info "enabling apache access logs"
    sed -i "s@#CustomLog logs/access.log common@CustomLog /dev/stdout common@" /etc/apache2/httpd.conf
fi

sed -i "s#%%servername%%#${servername}#" /etc/apache2/conf.d/monica.conf
sed -i "s#%%serveradmin%%#${serveradmin}#" /etc/apache2/conf.d/monica.conf
sed -i "s#%%certpath%%#${certpath}#g" /etc/apache2/conf.d/monica.conf
sed -i "s#%%keypath%%#${keypath}#g" /etc/apache2/conf.d/monica.conf

bashio::log.info "starting apache2"
exec httpd -DFOREGROUND
